# supervisord.conf
[supervisord]
nodaemon=true

[program:kea]
command=/startup.sh

[program:kea-ctrl-agent]
command=/startupca.sh
