# supervisord.conf
[supervisord]
nodaemon=true

[program:kea]
command=/startup.sh

[program:apache2]
command=/startupca.sh
